# Kubernetes/Deployment/systemd/container.slice
[Unit]
Description=Limited resources Slice for containers and pods
DefaultDependencies=no
Before=slices.target
Requires=-.slice
After=-.slice
[Slice]
MemoryLimit=4G