# Kubernetes/Deployment/systemd/container.slice
[Unit]
Description=Limited resources Slice
DefaultDependencies=no
Before=slices.target

[Slice]
MemoryLimit=48G