# Kubernetes/Deployment/systemd/container.slice
[Unit]
Description=Limited resources Slice
DefaultDependencies=no
Before=slices.target

[Slice]
MemoryHigh=48G
MemoryMax=52G