# Kubernetes/Deployment/monitoring/Out-Cluster/node-exporter/node-exporter.service
[Unit]
Description=Node Exporter

[Service]
Slice=container.slice
ExecStart=/usr/bin/node_exporter

[Install]
WantedBy=multi-user.target
